# ft_services/srcs/wordpress/srcs/supervisor.conf
[supervisord]
nodaemon=true
autorestart=true

[program:script]
command=sh /tmp/start.sh
startsecs=0
autorestart=false
startretries=1

[program:nginx]
command=service nginx start
